# cron.d/send_email_with_gifs
MAILTO=anatolyburtsev@yandex.ru
path="/root/nofity_about_news_on_site"
0	*/4		*		*		*		root	python ${path}/mail/send_files_to_mail.py
